# pi-hole/rootfs/etc/services.d/dnscrypt-proxy/finish
#!/bin/bash

set -ueo pipefail

killall dnscrypt-proxy
